# com.bonny.self-improvement.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bonny.self-improvement</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rohan.vadgaonkar/workspace/junior-pm-agent-bonny/scripts/self-improvement.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/rohan.vadgaonkar/workspace/junior-pm-agent-bonny</string>

    <key>StandardOutPath</key>
    <string>/Users/rohan.vadgaonkar/workspace/junior-pm-agent-bonny/logs/self-improvement-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rohan.vadgaonkar/workspace/junior-pm-agent-bonny/logs/self-improvement-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/rohan.vadgaonkar</string>
    </dict>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
